# uniproc-agent.service
[Unit]
Description=Uniproc telemetry agent for WSL2
After=network.target
StartLimitBurst=3

[Service]
ExecStart=/usr/local/bin/uniproc-agent
Restart=on-failure
RestartSec=5s

AmbientCapabilities=CAP_BPF CAP_NET_ADMIN CAP_PERFMON CAP_SYSLOG
CapabilityBoundingSet=CAP_BPF CAP_NET_ADMIN CAP_PERFMON CAP_SYSLOG
LimitMEMLOCK=infinity

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectClock=true
RestrictAddressFamilies=AF_VSOCK
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target